# linux/etc/systemd/system/reflector.timer
[Unit]
Description=Periodic pacman mirrorlist refresh

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
