# autostart/keybase_autostart.desktop
# This file is generated the first time Keybase starts, along with a sentinel
# file at ~/.config/keybase/autostart_created. As long as the sentinel exists,
# this file won't be automatically recreated, so you can edit it or delete it
# as you like. Manually changing autostart settings in the GUI can stomp on
# your edits here though.

[Desktop Entry]
Name=Keybase
Comment=Keybase Filesystem Service and GUI
Type=Application
Exec=env KEYBASE_AUTOSTART=1 keybase-gui
