# linux/autostart/blueman.desktop
[Desktop Entry]
Type=Application
Exec=/bin/false
X-GNOME-Autostart-enabled=false
Hidden=true
