# apps/systemd/user/dconf-save.service
[Unit]
Description=Save Gnome, Cinnamon et other GTK apps configuration

[Service]
ExecStart= /home/visi/Bin/dconf-save
