# apps/systemd/user/pacsave.service
[Unit]
Description=Backup list of installed packages

[Service]
ExecStart=/home/visi/bin/pacsave